# crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 8,14,22 * * * python /app/main.py >> /var/log/cron.log 2>&1
