# crontab
0 9,16,22 * * * root python /app/main.py >> /var/log/cron.log 2>&1
